# etc/grooved-notify.service
[Unit]
Description=Desktop notifications support for Grooved
Requires=grooved.service

[Service]
Type=simple
ExecStart=/usr/bin/grooved-notify
